# Systemd/akula_sensors.service
# akula_sensors.service

[Unit] 
Description=Akula sensors 
Documentation=

[Service] 
Type=simple 

ExecStart=/home/mrob/AkulaRecord.sh
